# deploy/systemd/mobile-server-tunnel.service
[Unit]
Description=ERPNext Mobile Server Cloudflare Tunnel
After=network-online.target mobile-server-core.service
Wants=network-online.target mobile-server-core.service
StartLimitIntervalSec=0

[Service]
Type=simple
User=wikki
WorkingDirectory=/home/wikki/deploy/mobile_server_deploy
ExecStartPre=/bin/bash -lc 'pkill -f "cloudflared.*accord-vision-core" || true; sleep 1; pkill -9 -f "cloudflared.*accord-vision-core" || true; rm -f /home/wikki/deploy/mobile_server_deploy/.core_domain_tunnel.pid'
ExecStart=/bin/bash -lc 'exec /usr/local/bin/cloudflared tunnel --protocol http2 --config /home/wikki/deploy/mobile_server_deploy/.core_domain_tunnel.yml run accord-vision-core >> /home/wikki/deploy/mobile_server_deploy/.core_domain_tunnel.log 2>&1'
Restart=always
RestartSec=3
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
